# bin/nolost.service
[Unit]
Description=Nolost daemon
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mem daemon --port 7456
Restart=on-failure
RestartSec=3
Environment=MEM_HOME=%h/.local/share/nolost
# Uncomment to surface Python parsers (Docling/MarkItDown):
# Environment=MEM_PYTHON=%h/.local/share/nolost/python/bin/python
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
